# recipes-go/3mdeb-rtectrl/3mdeb-rtectrl/RteCtrl.service
[Unit]
Description=RTE REST API server
After=syslog.target

[Service]
Type=forking
ExecStart=/usr/bin/RteCtrl
Restart=on-failure

[Install]
WantedBy=multi-user.target
